# supervisord.conf
[supervisord]
nodaemon=true

[program:racktables]
command=/usr/bin/php-fpm5

[program:nginx]
command=/usr/sbin/nginx
